# systemd/screen-auto.timer
[Unit]
Description=Poll Night Light ramp every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target
